# web/php7/supervisor/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm --nodaemonize --fpm-config /etc/php/fpm/php-fpm.conf
